# snapshots/ip-172-31-11-160/20251007T174731Z/systemd_sidekiq_template.service
# /etc/systemd/system/sidekiq-flagship@.service
[Unit]
Description=Sidekiq (flagship) instance %i
After=network.target
Wants=sidekiq-flagship.slice

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# rbenv 경로 보장 (bundle not found 방지)
Environment=RBENV_ROOT=/home/ubuntu/.rbenv
Environment=PATH=/home/ubuntu/.rbenv/shims:/home/ubuntu/.rbenv/bin:/usr/local/bin:/usr/bin:/bin

# 공통/인스턴스 환경파일 (앱 비밀, DB/Redis URL, 키 등)
EnvironmentFile=-/etc/flagship-seoulmake.env
EnvironmentFile=-/etc/default/flagship-seoulmake
EnvironmentFile=/etc/sidekiq/flagship.env
EnvironmentFile=-/etc/sidekiq/flagship.d/%i.env

WorkingDirectory=/home/ubuntu/flagship-seoulmake
ExecStart=/usr/local/bin/run-sidekiq-flagship.sh %i

Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sidekiq-flagship[%i]
Slice=sidekiq-flagship.slice
LimitNOFILE=65536

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full
ReadWritePaths=/home/ubuntu/flagship-seoulmake /tmp /var/tmp
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target

# /etc/systemd/system/sidekiq-flagship@.service.d/env.conf
[Service]
EnvironmentFile=/etc/flagship-seoulmake.env
